# deploy/systemd/nervos-talk-forum-ingest.service
[Unit]
Description=Nervos Brain Talk forum incremental ingest
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=PROJECT_ROOT=%h/path/to/Nervos-Brain
Environment=MAMBA_BIN=%h/miniforge3/bin/mamba
Environment=MAMBA_ENV=nervos-brain
Environment=TALK_LATEST_PAGES=3
ExecStart=/bin/bash -lc 'cd "$PROJECT_ROOT" && "$MAMBA_BIN" run -n "$MAMBA_ENV" python scripts/run_talk_forum_ingest.py --latest-pages "$TALK_LATEST_PAGES" --incremental'
