# deploy/systemd/nervos-github-ingest.service
[Unit]
Description=Nervos Brain GitHub docs/code incremental ingest
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=PROJECT_ROOT=%h/path/to/Nervos-Brain
Environment=MAMBA_BIN=%h/miniforge3/bin/mamba
Environment=MAMBA_ENV=nervos-brain
EnvironmentFile=-%h/.config/nervos-brain/github-ingest.env
ExecStart=/bin/bash -lc 'cd "$PROJECT_ROOT" && "$MAMBA_BIN" run -n "$MAMBA_ENV" python scripts/run_github_docs_ingest.py --incremental && "$MAMBA_BIN" run -n "$MAMBA_ENV" python scripts/run_github_code_ingest.py --incremental'
